# backend/gunicorn/gunicorn.service
[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target
[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/healthlinkdrf
ExecStart=/home/ubuntu/env/bin/gunicorn --access-logfile - --workers 3 --bind unix:/home/ubuntu/healthlinkdrf/healthlink/healthlink.sock healthlink.wsgi:application

       
[Install]
WantedBy=multi-user.target